# services/compute/trading-service/scripts/indicator-service.service
[Unit]
Description=Indicator Calculation Service
After=network.target postgresql.service

[Service]
Type=simple
WorkingDirectory=%h/tradecat/tradecat/services/compute/trading-service
Environment=PYTHONPATH=%h/tradecat/tradecat/services/compute/trading-service/src
ExecStart=/usr/bin/python3 -m indicator_service --daemon --log-file %h/tradecat/tradecat/services/compute/trading-service/logs/daemon.log
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
